# roles/n8n/templates/n8n.service.j2
[Unit]
Description=n8n workflow automation
After=network.target

[Service]
Type=simple
User={{ n8n.user }}
Group={{ n8n.group }}
WorkingDirectory={{ n8n.home }}
EnvironmentFile={{ n8n.data_dir }}/.env
ExecStart=/usr/local/bin/n8n start
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=n8n

# Hardening
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
PrivateDevices=true

[Install]
WantedBy=multi-user.target
